# backend/closecheck.service
[Unit]
Description=CloseCheck API
After=network.target

[Service]
Type=notify
User=app
Group=app
WorkingDirectory=/var/www/closecheck-backend
Environment="PATH=/var/www/closecheck-backend/venv/bin"
EnvironmentFile=/var/www/closecheck-backend/.env.prod
ExecStart=/var/www/closecheck-backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001 --workers 4
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=closecheck

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
